# sonar/sonar.supervisor.conf
[program:sonarqube]
command=/opt/sonarqube/bin/run.sh
stdout_logfile=/var/log/supervisor/%(program_name)s.log
stderr_logfile=/var/log/supervisor/%(program_name)s.log
autorestart=true
