# ubuntu-mysql/mysql.supervisor.conf
[program:mysql]
command=/entrypoint.sh mysqld
stdout_logfile=/var/log/supervisor/%(program_name)s.log
stderr_logfile=/var/log/supervisor/%(program_name)s.log
autorestart=true
