# tomcat/7/tomcat.supervisor.conf
[program:tomcat7]
command=/usr/local/tomcat/bin/catalina.sh run
autorestart=true
