# jenkins/jenkins.supervisor.conf
[program:jenkins]
command=/etc/init.d/jenkins start
stdout_logfile=/var/log/supervisor/%(program_name)s.log
stderr_logfile=/var/log/supervisor/%(program_name)s.log
autorestart=true
