# chaincoin/supervisord.conf
[supervisord]
nodaemon=true

[program:chaincoind]
command=chaincoind -daemon
stdout_logfile=/var/log/supervisor/%(program_name)s.log
stderr_logfile=/var/log/supervisor/%(program_name)s.log
